# tools/tools/nanobsd/gpt/Files/etc/rc.d/efiboot_init
#!/bin/sh
#
# Copyright (c) 2026 [XXX]
#
# SPDX-License-Identifier: BSD-2-Clause
#

# PROVIDE: efiboot_init
# REQUIRE: FILESYSTEMS
# BEFORE: LOGIN
# KEYWORD: firstboot

. /etc/rc.subr

name="efiboot_init"
desc="EFI BootOrder initializer"
start_cmd="efiboot_init_start"
start_postcmd="efiboot_init_poststart"
stop_cmd=":"
rcvar="efiboot_init_enable"

efiboot_init_start()
{
	local EFI_BOOT_VARS
	local alt_esp current_root efi_altlabel efi_label
	local id_1 id_2 nano_altroot nano_root order p_alternate p_current

	EFI_BOOT_VARS=$(efibootmgr)
	if echo "$EFI_BOOT_VARS" | grep -q "Primary ESP" && \
	    echo "$EFI_BOOT_VARS" | grep -q "Secondary ESP"; then
		startmsg "EFIBOOT INIT: NVRAM Boot Entries already provisioned."
		# XXXJL In our final version, we should remove existing ones,
		# and re-create them, as we don't knwow if the indices changed.
		return 0
	else
		startmsg "EFIBOOT INIT: NVRAM labels missing.  Executing First-Boot Configuration..."
	fi

	[ -f /etc/nanobsd.conf ] || err 1 "/etc/nanobsd.conf missing."
	. /etc/nanobsd.conf

	if [ -z "$NANO_DRIVE" ] || [ -z "$NANO_ROOT" ] || [ -z "$NANO_ALTROOT" ]; then
		err 1 "Configuration variables are missing."
	fi

	startmsg "Initializing UEFI boot entries:"

	# Detect EFI System Partitions
	p_efiboot1=$(gpart show -l "${NANO_DRIVE}" | awk '$4 == "efiboot1" {print $3}')
	p_efiboot2=$(gpart show -l "${NANO_DRIVE}" | awk '$4 == "efiboot2" {print $3}')
	if [ -z "$p_efiboot1" ] || [ -z "$p_efiboot2" ]; then
		err 1 "ERROR: Failed detect efiboot1/efiboot2 partition indices."
	fi

	# Setup UEFI boot labels
	current_root=$(df -n / | tail -n1 | cut -d' ' -f1)
	case "$current_root" in
	"/dev/gpt/${NANO_ALTROOT}")
		p_current="$p_efiboot2"
		p_alternate="$p_efiboot1"
		alt_esp="efiboot1"
		efi_label="Secondary ESP"
		efi_altlabel="Primary ESP"
		nano_root="$NANO_ALTROOT"
		nano_altroot="$NANO_ROOT"
		;;
	"/dev/gpt/${NANO_ROOT}")
		p_current="$p_efiboot1"
		p_alternate="$p_efiboot2"
		alt_esp="efiboot2"
		efi_label="Primary ESP"
		efi_altlabel="Secondary ESP"
		nano_root="$NANO_ROOT"
		nano_altroot="$NANO_ALTROOT"
		;;
	*)
		err 1 "Invalid root device '$current_root'."
		;;
	esac

	# Mount EFI System Partitions
	if ! mount | grep -q "on /boot/efi "; then
		if mount /boot/efi > /dev/null 2>&1; then
			efi_umount="/boot/efi"
		fi
	fi
	if ! mount | grep -q "on /mnt "; then
		if mount_msdosfs /dev/gpt/"${alt_esp}" /mnt > /dev/null 2>&1; then
			efi_umount="${efi_umount} /mnt"
		fi
	fi

	# Create the boot entries
	efibootmgr -c -l "/boot/efi/EFI/FreeBSD/loader.efi" -a \
	    -L "${efi_label}" \
	    -e "vfs.root.mountfrom=ufs:/dev/gpt/${nano_root}" > /dev/null

	efibootmgr -c -l "/mnt/EFI/FreeBSD/loader.efi" -a \
	    -L "${efi_altlabel}" \
	    -e "vfs.root.mountfrom=ufs:/dev/gpt/${nano_altroot}" > /dev/null

	# Get the 4-digit hex IDs
	id_1=$(efibootmgr | awk -v lbl="${efi_label}" '$0 ~ lbl {gsub(/[+]?Boot|\*/, "", $1); print $1}')
	id_2=$(efibootmgr | awk -v lbl="${efi_altlabel}" '$0 ~ lbl {gsub(/[+]?Boot|\*/, "", $1); print $1}')

	if [ -z "$id_1" ] || [ -z "$id_2" ]; then
		err 1 "ERROR: Failed to get UEFI Boot IDs."
	fi

	# Apply the BootOrder
	order="${id_1},${id_2}"
	efibootmgr -o "${order}" > /dev/null

	startmsg "BootOrder updated successfully: ${order}"
}

efiboot_init_poststart()
{
	if [ -n "$efi_umount" ]; then
		umount ${efi_umount} >/dev/null 2>&1
	fi
}

load_rc_config $name
efiboot_init_svcj="NO"
run_rc_command "$1"
